# deploy/openhub-gateway.service
[Unit]
Description=OpenHub Rust Gateway
After=network.target openhub-hub.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/openhub/gateway
ExecStart=/root/openhub/gateway/bin/gateway
Restart=always
RestartSec=3
Environment=RUST_LOG=info
Environment=DATABASE_URL=postgres://xinference:password@localhost:5432/openhub

[Install]
WantedBy=multi-user.target
