# deploy/openhub-hub.service
[Unit]
Description=OpenHub Node Service (Hub)
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/openhub
ExecStart=/usr/bin/npm run dev -w @openhub/hub
Restart=always
RestartSec=3
EnvironmentFile=-/root/.openhub.env
Environment=PORT=3000
Environment=DATABASE_URL=postgres://xinference:password@localhost:5432/openhub

[Install]
WantedBy=multi-user.target
