# deploy/openhub-gateway.service
[Unit]
Description=OpenHub Rust Gateway
After=network.target openhub-hub.service

[Service]
Type=simple
User=mac-m4
WorkingDirectory=/home/mac-m4/openhub/gateway
ExecStart=/home/mac-m4/openhub/gateway/target/release/gateway
Restart=always
RestartSec=3
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
