# deploy/openhub-hub.service
[Unit]
Description=OpenHub Node Service (Hub)
After=network.target postgresql.service

[Service]
Type=simple
User=mac-m4
WorkingDirectory=/home/mac-m4/openhub
ExecStart=/usr/bin/npm run dev -w @openhub/hub
Restart=always
RestartSec=3
Environment=NODE_ENV=production
Environment=DATABASE_URL=postgres://xinference:password@localhost:5432/openhub
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
